# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:api]
directory=/opt/api
command=/opt/api/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:worker]
directory=/opt/api
command=/opt/api/.venv/bin/celery -A app.core.celery_app.celery_app worker --loglevel=info --concurrency=%(ENV_CELERY_WORKER_CONCURRENCY)s
autorestart=true
priority=15
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; -----------------------------------------------------------------
; Tur 4 — rule-based-engine orchestration workers.
; Bu üç process listener.db üzerinde çalışan LangGraph runtime
; bileşenleri. uvicorn API process'i ile aynı container içinde
; supervisord tarafından yönetiliyor.
; -----------------------------------------------------------------

[program:rbe-listener]
directory=/opt/api
command=/opt/api/.venv/bin/python listener.py
environment=INTERNAL_SERVICE_IN_PROCESS=1
autorestart=true
priority=16
startsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:rbe-workflow]
directory=/opt/api
command=/opt/api/.venv/bin/python workflow_worker.py
environment=INTERNAL_SERVICE_IN_PROCESS=1
autorestart=true
priority=17
startsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:rbe-task]
directory=/opt/api
command=/opt/api/.venv/bin/python task_executor.py
environment=INTERNAL_SERVICE_IN_PROCESS=1
autorestart=true
priority=17
startsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:php-fpm]
command=/usr/sbin/php-fpm8.2 -F
autorestart=true
priority=18
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
